# backlog/paused/BL-359-always-on-operator-presence.yaml
id: BL-359
title: "The Operator keeps dropping — the channel is built, but nothing supervises the processes that own it"
milestone: M8
type: feature
status: todo
severity: medium
priority: 4
direction: human-requested
source: "backlog/INTAKE-always-on-operator.md — the human via Telegram SUP-2, 2026-07-13T15:16Z: 'I want operator up & always on. File a bug ticket as it keeps dropping'."
depends_on: []

notes: |
  THE ASK. A persistent, always-available Operator the human can talk to at any time — not one that
  vanishes between messages.

  ── DESIGN FORK: RESOLVED. The human picked the DURABLE CONVERSATION SURFACE (2026-07-14). ────────
  "Always on" means the standing Operator topic (BL-346) IS the presence: always answerable, backed
  by DISPOSABLE runs, no idle token burn. It does NOT mean a long-lived interactive session. A
  permanently-attended Opus session is REJECTED — it burns tokens while idle and holds the single
  Operator slot. `attend.sh` is therefore NOT the always-on mechanism, and this ticket must stop it
  being mistaken for one.

  ── THE TICKET'S OWN PREMISE WAS WRONG. RE-VERIFIED IN CODE AND ON THIS HOST, 2026-07-14. ─────────
  Do not build from the previous draft of these notes; it asserted the runtime was already
  supervised, and it is not. What is actually true:

  ALREADY BUILT — the conversation surface is essentially DONE. Do not rebuild it.
    - The standing topic exists and is durably bound: `.swarmforge/operator/telegram-topic-map.json`
      carries `"292":"OPERATOR"`. Written by `ensureOperatorTopic`
      (`extension/src/tools/telegram-front-desk-bot.ts:230-244`), called from `main()` BEFORE the
      poll loops start (:606), so it is idempotent across restarts and cannot race.
    - The round trip is code-complete with NO dark hop: inbound (topic 292 → subject OPERATOR →
      `postToBridge`) → `bridgeServer.ts:297` emits `TELEGRAM_TOPIC_MESSAGE` into
      `.swarmforge/operator/events.jsonl` → `operator_runtime.bb` dispatches a disposable run
      (`launch-operator!`, or BL-334's restricted second lane `launch-front-desk-operator!`) →
      `append-to-reply-outbox!` → SSE → `resolveReplyTopicId` resolves OPERATOR back to topic 292.
    - BL-358's untagged-gate routing into that topic is genuinely wired (a real production caller in
      `buildConciergeTickAdapters`), not dark.

  THE REAL GAP 1 — NOTHING IS SUPERVISED ON THIS HOST. This is the build.
    - Verified: `/etc/systemd/system` has no swarm unit; there are no user units. `operator_runtime.bb`
      (pid 110641) and `front_desk_supervisor.bb` (pid 110530) both have PPID 1 (`/init`) — orphaned
      `nohup` children of `start_operator_runtime.sh:33`. A crash, an OOM kill, a WSL restart or a
      reboot takes the whole conversation surface down and NOTHING brings it back.
    - The systemd units are RENDERED BUT NEVER INSTALLED. `swarmforge/deploy/generate_systemd_units.sh`
      has correct `operator` (`Restart=always`, `WantedBy=multi-user.target`) and `front-desk`
      branches. But the ONLY caller is `swarmforge/deploy/provision_secondary_host.sh` — a
      SECONDARY-host provisioner never run against this box — and even it installs only the swarm and
      operator units, never `--unit=front-desk` (grep: that flag has zero non-test, non-doc callers).
      BL-351 shipped that branch and is in `done/`, but a generated-and-never-installed unit is, in
      its own runbook's words, "exactly as dark as no unit at all" — the epic-runtime-wiring rule
      firing again.
    - So the deliverable is an installer that actually enables BOTH units on the PRIMARY/dogfood host,
      idempotently and re-runnably, giving restart-on-crash and autostart-on-boot. This is what makes
      "the channel never dies" true rather than merely intended.

  THE REAL GAP 2 — THE SINGLE-SLOT INVARIANT IS A FICTION, AND MUST BE MADE HONEST.
    - `operator-running?` (`operator_runtime.bb:278-284`) decides the slot by exactly two things: a
      tmux session named `operator` on the Operator's OWN socket, or a live pid in `operator.pid`.
    - `.swarmforge/operator/attend.sh` is a bare `exec claude …` (untracked, hand-written; grep finds
      it referenced NOWHERE in the repo except this ticket). It creates NEITHER of those. So an
      attended session does NOT register as the slot-holder — even though `operator_runtime.bb:911-917`
      states in a comment that "an attended (interactive) session holds it indefinitely by design",
      and BL-334 + BL-345's starvation alarm both rest on that premise.
    - Consequence, today: with a human sitting in an attended session, `should-launch-operator?` still
      fires and the runtime spawns a SECOND, unrestricted, tool-holding Operator beside him; and
      BL-334's second lane is simultaneously NOT eligible (its gate requires
      `full-operator-running?`, which reads false). The invariant everything else is built on does not
      hold.
    - Resolve it, per this ticket's standing constraint (an always-on Operator that silently disables
      health sweeps, respawns and stall nudges is a STRICTLY WORSE failure mode than one that drops,
      because it is invisible). Under the chosen design the presence no longer holds the slot at all —
      disposable runs acquire and release it — so recovery keeps working. But an interactive session
      must never be able to break the invariant SILENTLY: either it registers as a genuine slot-holder
      (so the runtime sees it and BL-334's front-desk lane engages), or `attend.sh` is retired. Do not
      leave the code asserting one thing and doing another.

  GAP 3 — NEVER EXERCISED LIVE. The reply outbox contains ZERO `"threadId":"OPERATOR"` records: the
  topic was bound 2026-07-13 and the OPERATOR round trip has apparently never run end to end.
  Code-complete is not proven. QA must drive it live (below); a green unit suite is not the evidence
  this ticket needs.

  SCOPE. One parcel, two halves that must land together (scenario 05 depends on 06 being true):
    - Supervision: `swarmforge/deploy/generate_systemd_units.sh` (already has both branches — reuse,
      do not duplicate) plus a real installer for THIS host covering the operator unit AND the
      front-desk unit. `swarmforge/scripts/start_operator_runtime.sh`'s `nohup` path is what is being
      replaced as the primary start mechanism.
    - Slot honesty: `operator-running?` and the attended-mode launch path in
      `swarmforge/scripts/operator_runtime.bb`.
  Out of scope: rebuilding the standing topic, the reply routing, or the BL-334 second lane — all
  three already exist and are verified above. Also out of scope: making the DISPOSABLE Operator
  long-lived; the human explicitly rejected that.

  ⚠ FILE OVERLAP — SEQUENCE, DO NOT RUN CONCURRENTLY WITH BL-354. Both tickets edit
  `swarmforge/scripts/operator_runtime.bb` (BL-354: `write-telegram-reply-context!` + the await
  lifecycle; BL-359: `operator-running?` + the attended launch path). Different functions, same file,
  so concurrent promotion means a merge conflict. `depends_on` is empty because neither BLOCKS the
  other logically — this is a scheduling constraint for the coordinator, per the orthogonality rule.

  IMPLEMENTATION NOTES FOR THE CODER.
    - Scenario 04 is a Scenario Outline. Per the engineering rule, its `<mishap>` column must be
      LOAD-BEARING: validate against an explicit KNOWN_VALUES lookup ("a crash", "a host reboot") and
      throw on anything else — never a binary `===`/ternary that lumps the unmatched value into an
      `else`, or gherkin-mutator's mutant survives to the hardener.
    - Proof strategy for scenario 04: "a crash" is directly testable — kill the process and assert it
      returns with no human. "a host reboot" is NOT rebooted in a test; the honest mechanical proof is
      that the unit is installed AND enabled (autostart) with `Restart=always`. Assert that, and do
      not fake a reboot.
    - Do not use real timers/sleeps to wait for a restart (engineering: the no-real-timers ban covers
      subprocess and real-daemon tests too). Drive the supervisor's backoff through an injected seam,
      as the existing handoffd wiring tests do for their tunables.

  IR-DRY REVIEW (done at spec time — do not relitigate). One residual `possible-synonym` finding:
  "the Operator presence is live" vs "the Operator presence is lost to <mishap>". These are literal
  opposites; the shared tokens are inherent. Deliberate, not drift. The file is final as written.

  E2E QA PROCEDURE (this ticket's evidence is LIVE, not a unit suite — see GAP 3).
    1. Address the Operator in its standing Telegram topic; assert an answer returns to THAT topic and
       that a `"threadId":"OPERATOR"` record appears in `telegram-reply-outbox.jsonl` (today: zero).
    2. Let the disposable run finish and exit; address it again; assert it still answers.
    3. `systemctl kill` / `kill -9` the operator runtime; assert it returns with no human action, and
       that the Operator is reachable again.
    4. Assert both units report `enabled` (`systemctl is-enabled`) — the reboot-survival proof.
    5. Close the terminal the swarm was started from; assert the Operator is still reachable.
    6. WITH the presence live, kill a role's pane and leave a handoff unattended; assert the swarm's
       own health sweep still detects and recovers them — the presence must never hold the slot that
       recovery needs.
    7. Start an interactive Operator session; assert the runtime SEES an Operator as running and does
       not spawn a second unrestricted one alongside it.

acceptance:
  feature: specs/features/BL-359-always-on-operator-presence.feature
  steps:
    - "The human can reach the Operator at any moment"
    - "The Operator presence does not vanish when a run finishes"
    - "The Operator presence survives the terminal it was started from"
    - "The Operator presence comes back by itself, with no human to restart it"
    - "An always-on Operator never suspends the swarm's own recovery"
    - "An interactive Operator session can never go unseen by the swarm"

human_approval: pending

mutation_cost: medium
